# deploy/writerace.service
[Unit]
Description=WriteRace server
After=network.target

[Service]
Type=simple
User=writerace
Group=writerace
ExecStart=/usr/local/bin/writerace
Environment=PORT=8080
Environment=ALLOWED_ORIGIN=https://writerace.io
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
